# sickbeard.conf
[program:sickbeard]
user = nobody
command=/opt/sickbeard/SickBeard.py --config /config --datadir /data