# sickbeard.conf
[program:sickbeard]
user = nobody
command = /usr/sbin/python2 /opt/sickbeard/SickBeard.py --config /config/config.ini --datadir /config