# sickbeard.conf
[program:sickbeard]
user = nobody
command=/usr/bin/env python2 /opt/sickbeard/SickBeard.py --quiet --config /config --datadir /data