# sickbeard.conf
[program:sickbeard]
user = nobody
command=python2 /opt/sickbeard/SickBeard.py --config /config --datadir /data